# crontab
#run python script every minutes
@daily python /app/app.py