# crontab
#run bash script every day
@daily /app/scripts/crawl.sh
